# crontab
*/5 * * * *   root    chown -R jenkins:jenkins /var/lib/jenkins
